# debian/postinst
#!/bin/sh
/usr/sbin/update-rc.d -f sinkhole-service defaults
/etc/init.d/sinkhole-service start
