# debian/postrm
#!/bin/sh
if [ "$1" = "remove" ]; then
	/usr/sbin/update-rc.d -f sinkhole-service remove
fi
